# supervisord.conf
[supervisord]
nodaemon=true

[program:ffmpeg]
command = /usr/bin/ffmpeg
autostart = true

[program:ffprobe]
command = /usr/bin/ffprobe
autostart = true
